# daemon/config/gnode-daemon.service
[Unit]
Description=gNode (Geodineum Service Daemon) - Capability-based Service Discovery
Documentation=https://github.com/geodineum/gNode
# gnode-daemon needs ValKey, but the canonical valkey-gnode.service (created
# by Geodineum/install.sh ensure_canonical_valkey_unit) exists ONLY where the
# ValKey lives — the constellation master. Worker nodes reach it over the VPN
# and have no such unit, so a hard Requires= leaves a worker's daemon
# un-restartable (systemd: "Unit valkey-gnode.service not found"). Wants= keeps
# master ordering/propagation and tolerates absence; the daemon's own reconnect
# loop covers ValKey being down or remote.
After=network.target valkey-gnode.service
Wants=valkey-gnode.service
PartOf=valkey-gnode.service

[Service]
Type=simple
User=gnode
Group=gnode

# Environment configuration — single-canonical-route model (Commit 0.1):
# The daemon binary calls ecosystem_config::load() in main() to read the
# disk-minimal bootstrap.env (3 strict keys, root:root 0644) plus the
# ValKey-resident tier (geodineum:bootstrap:*). systemd no longer loads
# bootstrap.env directly.
#
# daemon.env remains an EnvironmentFile for daemon-specific tuning and
# VALKEY_PASSWORD_FILE path (read by ExecStartPre/ExecStart wrappers).
# .env is the dev-override hook (optional).
EnvironmentFile=-/etc/geodineum/components/gnode-daemon/daemon.env
EnvironmentFile=-/opt/geodineum/gNode/.env

# Production working directory
WorkingDirectory=/opt/geodineum/gNode

# Start the daemon with production settings (using ACL authentication)
# Multi-node routing: --node-id is the node's UNIQUE identity (and consumer
#                     name in the shared group). It MUST differ per node — two
#                     daemons sharing one name don't load-balance, they contend
#                     over the same pending entries. Read from GNODE_NODE_ID in
#                     daemon.env (the installer sets it to the host's name on a
#                     constellation join); defaults to "master" for the
#                     master/standalone node.
#                     --node-type is this node.s EXPOSURE: a comma-separated set
#                     of work classes it will process (an entry matches if any
#                     member covers its _gh hint). Read from GNODE_NODE_TYPE in
#                     daemon.env; defaults to "general". Manage with:
#                       geodineum node expose <hint[,hint...]>
# Tuning parameters (dimensions, batch sizes, etc.) are read from .env GNODE_* vars
# Node config path fallback chain (auto-detected):
#   1. --node-config-dir flag (if provided)
#   2. GNODE_NODE_CONFIG_DIR env var
#   3. /etc/geodineum/components/gnode-daemon/nodes/ (centralized ecosystem)
#   4. daemon/config/nodes/ (legacy)
# NOTE: --site-id and --environment removed - daemon discovers sites and listens to ALL DTAP environments
# Verify credential file exists before starting
# VALKEY_PASSWORD_FILE comes from daemon.env (default: /etc/geodineum/credentials/valkey_daemon.password)
ExecStartPre=/bin/bash -c '\
  PASS_FILE="${VALKEY_PASSWORD_FILE:-/etc/geodineum/credentials/valkey_daemon.password}"; \
  if [ ! -f "$PASS_FILE" ]; then \
    echo "Error: ValKey daemon password file not found at $PASS_FILE"; \
    exit 1; \
  fi; \
  PASS_CONTENT="$(cat "$PASS_FILE")"; \
  if [ "$PASS_CONTENT" = "CHANGE_ME_TO_STRONG_PASSWORD" ] || [ -z "$PASS_CONTENT" ]; then \
    echo "Error: Password file contains placeholder or is empty. Run setup-gnode.sh first."; \
    exit 1; \
  fi'

# Start the daemon. Connection parameters (host/port/user/topology/stream/
# log_level) flow from ecosystem_config::load() into clap via `env=`
# attributes — no need to pass them as CLI args. Wrapper only reads the
# password file into GNODE_REDIS_AUTH (keeps the password out of argv).
ExecStart=/bin/bash -c '\
  PASS_FILE="${VALKEY_PASSWORD_FILE:-/etc/geodineum/credentials/valkey_daemon.password}"; \
  export GNODE_REDIS_AUTH="$(cat $PASS_FILE)"; \
  exec ${GNODE_DAEMON_BIN:-/opt/geodineum/gNode/daemon/target/release/gnode-daemon} \
  --node-id "${GNODE_NODE_ID:-master}" \
  --node-type "${GNODE_NODE_TYPE:-general}" \
  start'

# Health check via built-in status command (daemon self-loads connection params)
ExecReload=/bin/bash -c '\
  PASS_FILE="${VALKEY_PASSWORD_FILE:-/etc/geodineum/credentials/valkey_daemon.password}"; \
  GNODE_REDIS_AUTH="$(cat $PASS_FILE)" \
  ${GNODE_DAEMON_BIN:-/opt/geodineum/gNode/daemon/target/release/gnode-daemon} \
  status'

# Graceful shutdown (daemon self-loads connection params)
ExecStop=/bin/bash -c '\
  PASS_FILE="${VALKEY_PASSWORD_FILE:-/etc/geodineum/credentials/valkey_daemon.password}"; \
  GNODE_REDIS_AUTH="$(cat $PASS_FILE)" \
  ${GNODE_DAEMON_BIN:-/opt/geodineum/gNode/daemon/target/release/gnode-daemon} \
  stop'

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/geodineum/gNode/.gnode /opt/geodineum/gNode/daemon/target /opt/geodineum/gNode/logs /var/log/geodineum/gnode /etc/geodineum

# Resource limits
LimitNOFILE=65536
# Hard backstop on thread/task count. Steady state is ~25-40 threads (all
# tokio pools are bounded in code: main min(cores,4), fast-lane 2, compute 2,
# asset 2). 512 leaves wide headroom while guaranteeing the cgroup kills any
# runaway spawn instead of exhausting the host.
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gnode-daemon

[Install]
WantedBy=multi-user.target
